# scripts/agent-session-analytics.service
[Unit]
Description=Agent Session Analytics - MCP server for session log analysis
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONPATH=__PROJECT_DIR__/src
ExecStart=__VENV_PYTHON__ -m agent_session_analytics.server
Restart=always
RestartSec=5
StandardOutput=append:__HOME__/.claude/contrib/agent-session-analytics/agent-session-analytics.log
StandardError=append:__HOME__/.claude/contrib/agent-session-analytics/agent-session-analytics.err

[Install]
WantedBy=default.target
