# scripts/claude-session-analytics.service
[Unit]
Description=Claude Session Analytics - MCP server for session log analysis
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONPATH=__PROJECT_DIR__/src
ExecStart=__VENV_PYTHON__ -m session_analytics.server
Restart=always
RestartSec=5
StandardOutput=append:__HOME__/.claude/session-analytics.log
StandardError=append:__HOME__/.claude/session-analytics.err

[Install]
WantedBy=default.target
