# scripts/com.evansenter.agent-session-analytics-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evansenter.agent-session-analytics-push</string>

    <key>ProgramArguments</key>
    <array>
        <string>__CLI_PATH__</string>
        <string>push</string>
        <string>--url</string>
        <string>__REMOTE_URL__</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude/contrib/agent-session-analytics/push.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/contrib/agent-session-analytics/push.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
